# firmament/docker/hadoop/runit/log/run
#!/bin/sh
exec chpst -u hadoop svlogd -tt /var/log/hadoop
