# firmament/docker/nifi/runit/log/run
#!/bin/sh
exec chpst -u nifi svlogd -tt /var/log/nifi
