# firmament/docker/kibana/5.6/runit/log/run
#!/bin/sh
exec chpst -u kibana svlogd -tt /var/log/kibana
