# firmament/docker/newman-phusion/3.1.0/runit/log/run
#!/bin/sh
exec chpst -u newman svlogd -tt /var/log/newman
